# celery.conf
; ==================================
;  celery worker supervisor example
; ==================================

[program:farm_celery]
; Set full path to celery program if using virtualenv
command=/home/deployer/farm/venv/bin/celery worker -A gearfarm --loglevel=INFO

; Alternatively,
;command=celery --app=your_app.celery:app worker --loglevel=INFO -n worker.%%h
; Or run a script
;command=celery.sh

directory=/home/deployer/farm/current
user=deployer
numprocs=1
stdout_logfile=/home/deployer/supervisor_logs/farm_celery.log
stderr_logfile=/home/deployer/supervisor_logs/farm_celery.err
autostart=true
autorestart=true
startsecs=10
environment=DJANGO_SETTINGS="gearfarm.my_settings.production_settings"

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000